# phpmyadmin/supervisord-httpd.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/bin/bash -c "exec /usr/sbin/apachectl -DFOREGROUND"
numprocs=1
autostart=true
autorestart=true
